# airootfs/etc/systemd/system/systemd-zram-setup@.service.d/archiso.conf
[Unit]
ConditionPathExists=!/run/archiso

[Service]
ExecStart=/bin/true
ExecStartPost=/bin/true
