# airootfs/etc/systemd/system/systemd-zram-generator.service.d/archiso.conf
[Unit]
ConditionPathExists=!/run/archiso

[Service]
ExecStart=
ExecStart=/bin/true
