# airootfs/etc/systemd/system/mados-fb-resolution.service
[Unit]
Description=madOS Framebuffer Resolution Configuration
Documentation=man:fbset(1)
After=systemd-modules-load.service
After=systemd-udev-trigger.service
After=plymouth-read-write.service
Before=display-manager.service
ConditionKernelCommandLine=!plymouth.enable=0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mados-fb-resolution.sh
RemainAfterExit=yes
TimeoutSec=30

[Install]
WantedBy=multi-user.target
