# airootfs/etc/systemd/system/mados-audio-init.service
[Unit]
Description=madOS Audio Initialization - Unmute ALSA Controls
After=sound.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/mados-audio-init.sh

[Install]
WantedBy=multi-user.target
