# airootfs/etc/systemd/system/mados-network-bootstrap.service
[Unit]
Description=madOS NetworkManager Ethernet Bootstrap
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mados-network-bootstrap.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
