# airootfs/etc/systemd/system/mados-timezone.service
[Unit]
Description=madOS Timezone Auto-Detection
Documentation=https://github.com/madkoding/mad-os
After=systemd-remount-fs.service network-online.target
Wants=network-online.target
Before=time-set.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/mados-timezone-detect.sh
StandardOutput=journal
StandardError=journal
TimeoutSec=60

[Install]
WantedBy=multi-user.target
