# airootfs/etc/systemd/system/mados-audio-init.service
[Unit]
Description=madOS Audio Initialization - Unmute ALSA Controls
# Run after sound hardware is available
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service sound.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/mados-audio-init.sh

[Install]
WantedBy=multi-user.target
