# airootfs/etc/systemd/system/mados-post-install.service
[Unit]
Description=madOS Post-Install Package Manager
# Wait for graphical session to be fully active
After=display-manager.service greetd.service systemd-user-sessions.service
Requires=display-manager.service greetd.service
# Network for package downloads
Wants=network-online.service
After=network-online.service
# Only run for first boot (not mados live user)
ConditionPathExists=!/home/mados/.config/mados/post-install-done
# Ensure a real user home exists
ConditionDirectoryNotEmpty=/home

[Service]
Type=oneshot
# Root is required for pacman operations
User=root
# Find first real user and get their display auth
ExecStartPre=/bin/sh -c 'for d in /home/*/.Xauthority; do [ -f "$d" ] && export XAUTHORITY=$d && export DISPLAY=:0 && break; done'
# Run post-install as root but with user's display permissions
ExecStart=/bin/bash -c 'export DISPLAY=:0; export XAUTHORITY=$(find /home -name ".Xauthority" -type f 2>/dev/null | head -1); /usr/local/bin/mados-post-install'
# Allow time for X11/Wayland to be ready
TimeoutStartSec=180
# Log everything
StandardOutput=journal
StandardError=journal
# Keep user environment
PassEnvironment=DISPLAY XAUTHORITY

[Install]
WantedBy=multi-user.target
