# airootfs/etc/systemd/system/mados-hw-quirks.service
[Unit]
Description=madOS hardware quirks (hardware-conditional)
DefaultDependencies=no
After=systemd-modules-load.service
Before=NetworkManager.service
Wants=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mados-hw-quirks.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
